# backup/backup.timer
[Unit]
Description="Automatically run backups"

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
